# provisions/templates/supervisord/magazine.ini
[program:magazine]
command=/home/vagrant/.virtualenvs/moxy/bin/gunicorn wsgi.magazine:application -w 2 -b 127.0.0.1:8000
directory=/vagrant/src/bilna-django/moxy
user=vagrant
autostart=true
autorestart=true
redirect_stderr=true
