# provisions/templates/supervisord/sellercenter.ini
[program:sellercenter]
command=/home/vagrant/.virtualenvs/moxy/bin/gunicorn wsgi.sellercenter:application -w 2 -b 127.0.0.1:8001
directory=/vagrant/src/moxy/moxy
user=vagrant
autostart=true
autorestart=true
redirect_stderr=true
